# overseerr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

# Set variables
slug=overseerr
port=5055
CONFIG_LOCATION=/config/addons_config/"$slug"/settings.json

# Wait for transmission to become available
bashio::net.wait_for "$port" localhost 900

if bashio::config.has_value 'CONFIG_LOCATION'; then
  CONFIG_LOCATION="$(bashio::config 'CONFIG_LOCATION')"
  # Modify if it is a base directory
  if [[ "$CONFIG_LOCATION" == *.* ]]; then CONFIG_LOCATION="$(dirname $CONFIG_LOCATION)"; fi
  CONFIG_LOCATION="$CONFIG_LOCATION"/settings.json
fi

bashio::log.info "Starting NGinx..."
exec nginx
